# docs/supervisor-queue.conf
# /etc/supervisor/conf.d/impar-queue.conf
# Mantém o worker da fila (e-mails) rodando e reiniciando no boot.

[program:impar-queue]
process_name=%(program_name)s_%(process_num)02d
command=php8.4 /var/www/impar-enquetes/backend/artisan queue:work --tries=3 --timeout=90 --sleep=3
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/impar-enquetes/backend/storage/logs/queue.log
stopwaitsecs=90
