# docs/supervisor-reverb.conf
# /etc/supervisor/conf.d/impar-reverb.conf
# Mantém o servidor WebSocket (Reverb) rodando e reiniciando no boot.

[program:impar-reverb]
process_name=%(program_name)s
command=php8.4 /var/www/impar-enquetes/backend/artisan reverb:start
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/impar-enquetes/backend/storage/logs/reverb.log
stopwaitsecs=10
